# crontab
0 3 1 */2 * root sh /scripts/run_certbot.sh
