# crontab
0 */12 * * * root sh /scripts/run_certbot.sh
